# roles/organized/templates/organized-celery.service.j2
[Unit]
Description=Organized App Celery Worker
OnFailure=notify-failure@%n.service
After=network.target postgresql.service valkey-server.service
Requires=valkey-server.service

[Service]
Type=simple
User={{ organized_user | default('organized') }}
Group={{ organized_user | default('organized') }}
WorkingDirectory={{ organized_dir }}
EnvironmentFile=/etc/iiab/organized.env
# Machine-local overrides — never written by Ansible; loaded last, so it wins.
EnvironmentFile=-/etc/iiab/organized.env.local
ExecStart={{ organized_dir }}/.venv/bin/celery \
    -A organized_backend worker \
    --loglevel=info \
    --concurrency=1 \
    --logfile=/var/log/organized/celery.log \
    --pidfile=/run/organized/celery.pid
Restart=on-failure
RestartSec=10

# Systemd hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ organized_dir }} /var/log/organized /run/organized
NoNewPrivileges=yes
PrivateTmp=yes

RuntimeDirectory=organized

[Install]
WantedBy=multi-user.target
