# roles/ai-gateway/templates/ai-gateway-spiritual.timer.j2
[Unit]
Description=Timer for the morning spiritual nudge

[Timer]
OnCalendar=*-*-* 06:50:00 America/Port-au-Prince
Persistent=true

[Install]
WantedBy=timers.target
