# roles/ai-gateway/templates/ai-gateway-reminders.timer.j2
[Unit]
Description=Timer for firing due operator reminders (every 2 min)

[Timer]
OnCalendar=*:0/2:00
Persistent=true

[Install]
WantedBy=timers.target
