# roles/ai-gateway/templates/ai-gateway.service.j2
[Unit]
Description=AI Gateway Service
Documentation=https://github.com/deldesir/gateway
After=network.target litellm.service
Wants=litellm.service

[Service]
Type=simple
User={{ ai_gateway_user }}
WorkingDirectory={{ ai_gateway_dir }}
StandardOutput=journal
StandardError=journal

# Secrets loaded from restricted file (mode 0600, owned by service user)
EnvironmentFile={{ ai_gateway_dir }}/.env

# Non-sensitive config — managed by Ansible, edit via local_vars.yml
Environment="RIVEBOT_URL={{ ai_gateway_rivebot_url | default('http://localhost:8087') }}"
Environment="GATEWAY_PUBLIC_URL={{ ai_gateway_public_url | default('http://localhost:8086') }}"
Environment="SIYUAN_PUBLIC_URL={{ ai_gateway_siyuan_public_url | default(ai_gateway_siyuan_url) }}"
{% if ai_gateway_siyuan_notebook_id | default('') %}
Environment="SIYUAN_NOTEBOOK_ID={{ ai_gateway_siyuan_notebook_id }}"
{% endif %}
Environment="DEFAULT_PERSONA={{ ai_gateway_default_persona | default('konex-support') }}"

ExecStart={{ ai_gateway_uv_bin }} run uvicorn {{ ai_gateway_app }} \
    --host 0.0.0.0 \
    --port {{ ai_gateway_port }}

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
