# roles/ai-gateway/templates/ai-gateway-kgextract.timer.j2
[Unit]
Description=Timer for nightly drawer→KG extraction

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true

[Install]
WantedBy=timers.target
