# roles/ai-gateway/templates/ai-gateway-kgwiki.timer.j2
[Unit]
Description=Timer for nightly KG→Wiki render

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
