# roles/ai-gateway/templates/ai-gateway-radar.timer.j2
[Unit]
Description=Timer for the morning cockpit radar (reality scan → SiYuan + digest)

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
