# roles/ai-gateway/templates/ai-gateway-eveningwrap.timer.j2
[Unit]
Description=Timer for the evening wrap (~20:30, closes the day, seeds tomorrow)

[Timer]
OnCalendar=*-*-* 20:30:00 America/Port-au-Prince
Persistent=true

[Install]
WantedBy=timers.target
