# roles/ai-gateway/templates/ai-gateway-talksync.service.j2
[Unit]
Description=Nightly Organized→talkmaster assignment sync (auto-create talk revisions)
OnFailure=notify-failure@%n.service
Documentation=https://github.com/deldesir/gateway
After=network.target organized.service siyuan.service
Wants=organized.service

[Service]
Type=oneshot
User=iiab-admin
WorkingDirectory={{ ai_gateway_dir }}
EnvironmentFile={{ ai_gateway_dir }}/.env
# Machine-local overrides — never written by Ansible; loaded last, so it wins.
EnvironmentFile=-{{ ai_gateway_dir }}/.env.local
ExecStart={{ ai_gateway_dir }}/.venv/bin/python -m app.jobs.talk_sync
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
