# roles/ai-gateway/templates/ai-gateway-personsync.timer.j2
[Unit]
Description=Timer for nightly person-hub sync (Organized → RapidPro + SiYuan)

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
