# roles/ai-gateway/templates/ai-gateway-persontouches.timer.j2
[Unit]
Description=Timer for person touches (daily, after the nightly person sync)

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
