# roles/ai-gateway/templates/second-brain-backup.timer.j2
[Unit]
Description=Timer for nightly second-brain backup

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
