# roles/siyuan/templates/siyuan.service.j2
[Unit]
Description=SiYuan Note-Taking Server
OnFailure=notify-failure@%n.service
After=network.target

[Service]
Type=simple
User={{ siyuan_user }}
WorkingDirectory={{ siyuan_dir }}
Environment=RUN_IN_CONTAINER=true
# Secret loaded from restricted file (mode 0600, owned by service user)
EnvironmentFile={{ siyuan_dir }}/.env
{% if siyuan_version.split('.')[0] | int > 3 or (siyuan_version.split('.')[0] | int == 3 and siyuan_version.split('.')[1] | int >= 7) %}
# SiYuan >= 3.7.0 moved the kernel to a cobra CLI: serving now requires the
# explicit `serve` subcommand. accessAuthCode is read from the flag here.
ExecStart={{ siyuan_dir }}/kernel/SiYuan-Kernel serve \
  --workspace={{ siyuan_workspace }} \
  --port={{ siyuan_port }} \
  --mode=prod \
  --accessAuthCode=${SIYUAN_ACCESS_AUTH_CODE}
{% else %}
ExecStart={{ siyuan_dir }}/kernel/SiYuan-Kernel \
  --workspace={{ siyuan_workspace }} \
  --port={{ siyuan_port }} \
  --mode=prod
{% endif %}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
