# roles/rivebot/templates/rivebot.service.j2
[Unit]
Description=RiveBot Service (Brain Engine)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{ rivebot_home }}
Environment="RIVEBOT_PORT={{ rivebot_port }}"
Environment="RIVEBOT_BRAINS_DIR={{ rivebot_brains_dir }}"
Environment="SIYUAN_DATA_DIR={{ rivebot_siyuan_data_dir }}"
Environment="RIVEBOT_GATEWAY_URL={{ rivebot_gateway_url }}"
Environment="RIVEBOT_MACRO_TIMEOUT_S={{ rivebot_macro_timeout_s }}"

ExecStart=/root/.local/bin/uv run uvicorn rivebot.main:app --host 127.0.0.1 --port {{ rivebot_port }} --log-level info
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
