# roles/rapidpro/templates/rapidpro-flows-export.timer.j2
# Managed by Ansible (rapidpro role)
[Unit]
Description=Nightly RapidPro flow export (diffable history)

[Timer]
OnCalendar=*-*-* 04:10:00
Persistent=true

[Install]
WantedBy=timers.target
