# roles/ai_middleware/templates/ai-middleware.service.j2
[Unit]
Description=Konex AI Middleware Service
After=network.target postgresql.service

[Service]
User={{ ai_middleware_user }}
Group={{ ai_middleware_user }}
WorkingDirectory={{ ai_middleware_dir }}
EnvironmentFile={{ ai_middleware_dir }}/.env
ExecStart={{ ai_middleware_dir }}/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port {{ ai_middleware_port }} --log-level info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
